# services/mylar/run
#!/bin/bash
exec /sbin/setuser abc python /app/mylar/Mylar.py --nolaunch --datadir /config/mylar
