# root/etc/services.d/mylar/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /app/mylar/Mylar.py --quiet --nolaunch \
	--datadir /config/mylar
